# anacrontab
# /etc/anacrontab: configuration file for anacron

# See anacron(8) and anacrontab(5) for details.

SHELL=/bin/bash
HOME=/home/sidd
LOGNAME=root

# These replace cron's entries
1	5	cron.daily	run-parts --report /etc/cron.daily
7	10	cron.weekly	run-parts --report /etc/cron.weekly
@monthly	15	cron.monthly	run-parts --report /etc/cron.monthly
7	1	backup.weekly	/bin/bash /home/sidd/dotfiles/backup.sh > /home/sidd/dotfiles/backup.log 2>&1
